# api/supervisord-boulderapi.conf
[program:boulderapi]
command=/api/bin/console server:start 0.0.0.0:8000
startsecs=0
